# packages/longport/systemd/openalice.service
[Unit]
Description=OpenAlice + Longbridge MCP
After=network.target

[Service]
Type=simple
WorkingDirectory={{OPENALICE_ROOT}}
ExecStart=/bin/bash -c '/usr/bin/node {{OPENALICE_ROOT}}/dist/main.js & /usr/bin/node {{LONGBRIDGE_MCP_ROOT}}/index.js'
Restart=always
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
